# crontab
0 * * * * * python /var/local/copernicus/manage.py generate_special_reports;
